# 1.0/api_server-1.0.5/scripts/etc/supervisord-request_routers.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=/var/log/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=info                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/supervisord/               ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

; This is where you run individual Tornado instances.
; We run four; one per processor core.
; In development, we ran as many as four per core with no issues.
; If you're looking to minimize cpu load, run fewer processes.
; BTW, Tornado processes are single threaded.
; To take advantage of multiple cores, you'll need multiple processes.

[program:request_router-8000]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8000
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8001]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8001
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8002]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8002
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8003]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8003
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8004]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8004
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8005]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8005
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8006]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8006
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8007]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8007
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8008]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8008
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8009]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8009
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8010]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8010
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8011]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8011
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8012]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8012
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8013]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8013
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8014]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8014
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8015]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8015
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8016]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8016
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8017]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8017
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8018]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8018
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8019]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8019
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8020]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8020
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8021]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8021
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8022]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8022
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8023]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8023
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8024]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8024
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8025]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8025
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8026]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8026
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8027]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8027
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8028]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8028
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8029]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8029
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8030]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8030
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8031]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8031
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8032]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8032
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8033]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8033
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8034]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8034
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

[program:request_router-8036]
command=/usr/local/bin/python /usr/local/bin/request_router/controllers/request_router.pyc --port=8036
process_name=%(program_name)s ; process_name expr (default %(program_name)s)

